# automation/systemd-scripts/pg_backup.service
[Unit]
Description=PostgreSQL Backup Service
After=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/bash -c 'test -f $HOME/pg_backup.py || { echo "Backup script not found"; exit 1; } && chmod +x $HOME/pg_backup.py'
ExecStart=/usr/bin/python3 $HOME/pg_backup.py
User=$USER

[Install]
WantedBy=multi-user.target